# templates/mc-kiosk.service
[Unit]
Description=Mission Control + Grafana Dual Kiosk Display
After=default.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStartPre=/bin/bash -c 'killall -9 chromium 2>/dev/null; sleep 2; true'
ExecStart=/bin/bash -c '\
  chromium --app=http://heavy:3000 \
    --class=chromium-mc \
    --user-data-dir=/tmp/chromium-mc \
    --no-first-run --disable-infobars --disable-session-crashed-bubble \
    --disable-translate --disable-features=TranslateUI \
    --disable-component-update --disable-background-networking \
    --disable-gpu --disable-gpu-compositing --disable-smooth-scrolling \
    --enable-low-end-device-mode \
    --password-store=basic & \
  sleep 3; \
  chromium --app=http://localhost:3001 \
    --class=chromium-grafana \
    --user-data-dir=/tmp/chromium-grafana \
    --no-first-run --disable-infobars --disable-session-crashed-bubble \
    --disable-translate --disable-features=TranslateUI \
    --disable-component-update --disable-background-networking \
    --disable-gpu --disable-gpu-compositing --disable-smooth-scrolling \
    --enable-low-end-device-mode \
    --password-store=basic & \
  wait'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
