# templates/openclaw-stats-collector.timer.j2
[Unit]
Description=Collect OpenClaw node stats every 30 seconds

[Timer]
OnBootSec=15
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
