# templates/openclaw-router-api.service.j2
[Unit]
Description=OpenClaw Cluster Service (FastAPI)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=enrico
EnvironmentFile=/home/enrico/pi-cluster/scripts/.env.cluster
# Dedicated venv avoids constraint drift from user-site pip installs
# (fastapi 0.115.6 pins starlette<0.42; user site had starlette 1.0.0).
# Create / refresh with: scripts/install-openclaw-router-venv.sh
ExecStart=/home/enrico/.local/openclaw-router-venv/bin/python /home/enrico/pi-cluster/scripts/openclaw-cluster-service.py
Restart=always
RestartSec=5
TimeoutStopSec=30
MemoryMax=512M
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=HOME=/home/enrico

[Install]
WantedBy=multi-user.target
