# templates/openclaw-node-agent.service.j2
[Unit]
Description=OpenClaw Node Stats Agent
After=network-online.target

[Service]
Type=oneshot
User=enrico
Environment=OPENCLAW_NODE_NAME={{ openclaw_node_name }}
Environment=OPENCLAW_MASTER_URL=http://{{ openclaw_gateway_host }}:8520
ExecStart=/usr/bin/python3 /opt/workspace/scripts/openclaw-node-agent.py
