# templates/openclaw-node-agent.timer.j2
[Unit]
Description=Push node stats to master every 30 seconds

[Timer]
OnBootSec=20
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
