# templates/ssh-watchdog.timer.j2
[Unit]
Description=Run SSH watchdog every 2 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=15

[Install]
WantedBy=timers.target
