# templates/zeroclaw-node.service.j2
[Unit]
Description=ZeroClaw Node Host ({{ openclaw_node_name }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=enrico
ExecStart=/usr/local/bin/zeroclaw node run --host {{ openclaw_gateway_host }} --port {{ openclaw_gateway_port }} --display-name "{{ openclaw_node_name }}"
Restart=always
RestartSec=10
Environment=HOME=/home/enrico

[Install]
WantedBy=multi-user.target
